# luci-app-vnt2/root/etc/uci-defaults/99_luci-app-vnt2
#!/bin/sh

# 设置启动脚本权限
chmod 755 /etc/init.d/vnt2
# 启用开机自启
/etc/init.d/vnt2 enable

exit 0
