# luci-app-vnt2/root/etc/uci-defaults/luci-vnt2
#!/bin/sh

# 使用uci batch模式批量执行配置命令，-q静默模式（忽略错误），>/dev/null屏蔽输出
uci -q batch <<-EOF >/dev/null
    # 删除ucitrack中已有的最后一个vnt相关配置（避免重复）
    delete ucitrack.@vnt[-1]
    # 新增一个ucitrack配置项，类型为vnt
    add ucitrack vnt
    # 设置该ucitrack项的init参数为vnt（关联VNT的启动脚本/etc/init.d/vnt）
    set ucitrack.@vnt[-1].init=vnt
    # 提交ucitrack配置更改（写入配置文件）
    commit ucitrack
EOF

# 赋予VNT服务启动脚本可执行权限（确保能启动/停止服务）
chmod +x /etc/init.d/vnt

# 删除LuCI的索引缓存和模块缓存（避免LuCI界面不显示VNT配置项）
rm -rf /tmp/luci-indexcache /tmp/luci-modulecache

# 脚本正常退出
exit 0
