# luci-app-vnt2/root/etc/uci-defaults/luci-vnt2
#!/bin/sh

# V2 版本 UCI 跟踪配置初始化（适配 vnt2）
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@vnt2[-1]  # 移除旧的 vnt2 跟踪项（避免重复）
	add ucitrack vnt2         # 添加 vnt2 跟踪项
	set ucitrack.@vnt2[-1].init=vnt2  # 关联 V2 启动脚本
	commit ucitrack           # 提交配置
EOF

# 赋予 V2 启动脚本执行权限
chmod +x /etc/init.d/vnt2

# 清理 LuCI 缓存（确保页面实时生效）
rm -rf /tmp/luci-indexcache /tmp/luci-modulecache

exit 0
