# luci-app-vnt2/root/etc/uci-defaults/luci-vnt2

#!/bin/sh
# Copyright (C) 2024 OpenWrt.org

rm -f /tmp/luci-indexcache
exit 0
