# root/etc/uci-defaults/luci-vnt2
#!/bin/sh
# uci-defaults for luci-app-vnt2

# Create default config if not exists
[ -f /etc/config/vnt2 ] || {
    touch /etc/config/vnt2
}

exit 0
